# cryoprotect-postgres-health.service
[Unit]
Description=CryoProtect PostgreSQL Health Check
After=cryoprotect-postgres.service
Requires=cryoprotect-postgres.service
Documentation=https://github.com/yourusername/cryoprotect

[Service]
Type=oneshot
User=mushu
WorkingDirectory=/home/mushu/Projects/CryoProtect
ExecStart=/home/mushu/Projects/CryoProtect/check_postgres_health.py
Environment=PATH=/home/mushu/.local/bin:/home/mushu/.npm-global/bin:/home/mushu/.local/bin:/usr/local/bin:/usr/bin:/var/lib/snapd/snap/bin:/usr/local/sbin:/usr/sbin

[Install]
WantedBy=multi-user.target