# cryoprotect-db-integrity.service
[Unit]
Description=CryoProtect Database Integrity Verification
After=network.target

[Service]
Type=oneshot
User=mushu
WorkingDirectory=/home/mushu/Projects/CryoProtect
ExecStart=/home/mushu/Projects/CryoProtect/run_verify_all_data_integrity.sh
Environment=PATH=/home/mushu/.local/bin:/home/mushu/.npm-global/bin:/home/mushu/.local/bin:/usr/local/bin:/usr/bin:/var/lib/snapd/snap/bin:/usr/local/sbin:/usr/sbin

[Install]
WantedBy=multi-user.target